# plugins/plex/extdrive.rules
ACTION=="add", SUBSYSTEM=="block", RUN+="/bin/sh -c '/usr/src/scripts/mount.sh'"
ACTION=="remove", SUBSYSTEM=="block", RUN+="/bin/sh -c '/usr/src/scripts/unmount.sh'"
